# services/pgbackrest-log/run
#!/bin/sh
exec logutil-service /var/log/pgbackrest-log
