# services/pgbackrest/run
#!/command/execlineb -P

with-contenv

/usr/bin/pgbackrest server
