# services/postgres/run
#!/command/execlineb -P
with-contenv
/usr/local/bin/docker-entrypoint.sh postgres
